# lat5150drvmil/99-archive/deployment-backup/debian-packages/dell-milspec-tpm2-dkms/DEBIAN/prerm
#!/bin/bash
set -e

case "$1" in
    remove|upgrade|deconfigure)
        systemctl stop tpm2-acceleration-early.service 2>/dev/null || true
        systemctl disable tpm2-acceleration-early.service 2>/dev/null || true
        modprobe -r tpm2_accel_early 2>/dev/null || true
        if [ "$1" = "remove" ]; then
            dkms remove -m tpm2-accel-early -v 1.0.0 --all 2>/dev/null || true
        fi
        ;;
esac

exit 0
