# lat5150drvmil/99-archive/root-cleanup-backup-20251011-165010/files/thermal-guardian.service
[Unit]
Description=Thermal Guardian - Advanced Thermal Management System
Documentation=file:///opt/github/LAT5150DRVMIL/README.md
After=multi-user.target
Wants=multi-user.target

[Service]
Type=simple
User=root
Group=root
ExecStart=/usr/bin/python3 /opt/github/LAT5150DRVMIL/thermal_guardian.py --daemon
ExecReload=/bin/kill -HUP $MAINPID
ExecStop=/bin/kill -TERM $MAINPID
KillMode=mixed
KillSignal=SIGTERM
TimeoutStopSec=30
Restart=always
RestartSec=5

# Security settings
NoNewPrivileges=false
PrivateTmp=true
ProtectSystem=false
ProtectHome=true
ReadWritePaths=/var/log /sys/class/hwmon /sys/class/thermal /sys/devices/system/cpu

# Resource limits
LimitNOFILE=1024
LimitNPROC=32

# Environment
Environment=PYTHONUNBUFFERED=1

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=thermal-guardian

[Install]
WantedBy=multi-user.target