# lat5150drvmil/packaging/debian/DEBIAN/postinst
#!/bin/bash
# postinst script for tpm2-accel-early-dkms
# Copyright (C) 2025 Military TPM2 Acceleration Project
# Classification: UNCLASSIFIED // FOR OFFICIAL USE ONLY

set -e

PACKAGE_NAME="tpm2-accel-early"
PACKAGE_VERSION="1.0.0"
MODULE_NAME="tpm2_accel_early"

case "$1" in
    configure)
        echo "Configuring tpm2-accel-early-dkms..."

        # Add module to DKMS
        echo "Adding module to DKMS..."
        dkms add -m "${PACKAGE_NAME}" -v "${PACKAGE_VERSION}" || true

        # Build module for current kernel
        echo "Building module for kernel $(uname -r)..."
        dkms build -m "${PACKAGE_NAME}" -v "${PACKAGE_VERSION}" -k "$(uname -r)" || {
            echo "ERROR: DKMS build failed" >&2
            exit 1
        }

        # Install the built module
        echo "Installing module..."
        dkms install -m "${PACKAGE_NAME}" -v "${PACKAGE_VERSION}" -k "$(uname -r)" || {
            echo "ERROR: DKMS install failed" >&2
            exit 1
        }

        # Create modprobe configuration
        echo "Creating modprobe configuration..."
        cat > /etc/modprobe.d/tpm2-acceleration.conf <<'EOF'
# TPM2 Early Boot Acceleration Module Configuration
# Auto-generated by tpm2-accel-early-dkms package

# Module parameters
# security_level: 0=UNCLASSIFIED, 1=CONFIDENTIAL, 2=SECRET, 3=TOP_SECRET
# early_init: 1=enable early boot initialization
# debug_mode: 0=disabled, 1=enabled

options tpm2_accel_early early_init=1 debug_mode=0 security_level=0

# Prevent automatic loading (systemd service handles this)
install tpm2_accel_early /bin/true
EOF

        # Create modules-load configuration for early boot
        echo "Creating modules-load configuration..."
        cat > /etc/modules-load.d/tpm2-acceleration.conf <<'EOF'
# TPM2 Early Boot Acceleration Module
# This ensures the module is loaded during early boot
# Auto-generated by tpm2-accel-early-dkms package

tpm2_accel_early
EOF

        # Update initramfs to include module in early boot
        echo "Updating initramfs for early boot support..."
        if command -v update-initramfs >/dev/null 2>&1; then
            update-initramfs -u -k all || {
                echo "WARNING: Failed to update initramfs" >&2
            }
        fi

        # Enable and start systemd service
        if [ -d /run/systemd/system ]; then
            echo "Enabling TPM2 acceleration service..."
            systemctl daemon-reload || true
            systemctl enable tpm2-acceleration-early.service || true

            # Load module immediately if TPM hardware is available
            if [ -e /dev/tpm0 ] || [ -e /dev/tpmrm0 ]; then
                echo "Loading module for current session..."
                modprobe "${MODULE_NAME}" early_init=1 security_level=0 || {
                    echo "WARNING: Failed to load module (may require reboot)" >&2
                }
            else
                echo "NOTE: TPM hardware not detected, module will load on next boot"
            fi
        fi

        # Create device node if it doesn't exist (usually handled by udev)
        if [ ! -e "/dev/${MODULE_NAME}" ]; then
            echo "NOTE: Device /dev/${MODULE_NAME} will be created when module loads"
        fi

        # Display status information
        echo ""
        echo "========================================="
        echo "TPM2 Early Boot Acceleration Installed"
        echo "========================================="
        echo "Package: ${PACKAGE_NAME}-dkms"
        echo "Version: ${PACKAGE_VERSION}"
        echo "Module:  ${MODULE_NAME}"
        echo "Device:  /dev/${MODULE_NAME}"
        echo ""
        echo "Configuration files:"
        echo "  /etc/modprobe.d/tpm2-acceleration.conf"
        echo "  /etc/modules-load.d/tpm2-acceleration.conf"
        echo ""
        echo "Default security level: UNCLASSIFIED (0)"
        echo ""
        echo "To change security level, edit:"
        echo "  /etc/modprobe.d/tpm2-acceleration.conf"
        echo "  and set: options tpm2_accel_early security_level=<0-3>"
        echo ""
        echo "Then run: update-initramfs -u && reboot"
        echo "========================================="
        echo ""
        ;;

    abort-upgrade|abort-remove|abort-deconfigure)
        ;;

    *)
        echo "postinst called with unknown argument \`$1'" >&2
        exit 1
        ;;
esac

# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.

#DEBHELPER#

exit 0
