# lat5150drvmil/packaging/debian/DEBIAN/prerm
#!/bin/bash
# prerm script for tpm2-accel-early-dkms
# Copyright (C) 2025 Military TPM2 Acceleration Project
# Classification: UNCLASSIFIED // FOR OFFICIAL USE ONLY

set -e

PACKAGE_NAME="tpm2-accel-early"
PACKAGE_VERSION="1.0.0"
MODULE_NAME="tpm2_accel_early"

case "$1" in
    remove|upgrade|deconfigure)
        echo "Preparing to remove tpm2-accel-early-dkms..."

        # Stop systemd service if running
        if [ -d /run/systemd/system ]; then
            echo "Stopping TPM2 acceleration service..."
            systemctl stop tpm2-acceleration-early.service 2>/dev/null || true
            systemctl disable tpm2-acceleration-early.service 2>/dev/null || true
        fi

        # Unload the module if it's loaded
        echo "Unloading module ${MODULE_NAME}..."
        if lsmod | grep -q "^${MODULE_NAME} "; then
            rmmod "${MODULE_NAME}" 2>/dev/null || {
                echo "WARNING: Could not unload ${MODULE_NAME} (may be in use)" >&2
                echo "A reboot may be required to complete removal" >&2
            }
        fi

        # Remove from DKMS only on actual removal (not upgrade)
        if [ "$1" = "remove" ]; then
            echo "Removing module from DKMS..."

            # Get all installed kernels for this module
            INSTALLED_KERNELS=$(dkms status -m "${PACKAGE_NAME}" -v "${PACKAGE_VERSION}" 2>/dev/null | \
                               awk -F', ' '{print $2}' | sed 's/kernel //' || true)

            # Uninstall from each kernel
            for KERNEL in ${INSTALLED_KERNELS}; do
                echo "Uninstalling from kernel ${KERNEL}..."
                dkms uninstall -m "${PACKAGE_NAME}" -v "${PACKAGE_VERSION}" -k "${KERNEL}" 2>/dev/null || true
            done

            # Remove DKMS source
            dkms remove -m "${PACKAGE_NAME}" -v "${PACKAGE_VERSION}" --all 2>/dev/null || true
        fi
        ;;

    failed-upgrade)
        ;;

    *)
        echo "prerm called with unknown argument \`$1'" >&2
        exit 1
        ;;
esac

# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.

#DEBHELPER#

exit 0
