# lat5150drvmil/99-archive/deployment-backup/debian-packages/dell-milspec-tpm2-dkms/lib/systemd/system/tpm2-acceleration-early.service
[Unit]
Description=TPM2 Hardware Acceleration Early Boot Integration
Documentation=man:tpm2_accel_early(8)
After=multi-user.target
Wants=multi-user.target
DefaultDependencies=no

[Service]
Type=oneshot
ExecStart=/bin/bash -c "test -c /dev/tpm2_accel_early && echo 'TPM2 acceleration ready' || echo 'TPM2 acceleration pending'"
RemainAfterExit=yes
StandardOutput=journal
StandardError=journal
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=yes
NoNewPrivileges=yes

[Install]
WantedBy=multi-user.target
