# lat5150drvmil/deployment/systemd/lat5150-self-improvement.service
[Unit]
Description=LAT5150 AI Self-Improvement Service
Documentation=file:///home/user/LAT5150DRVMIL/02-ai-engine/ai_self_improvement.py
After=network-online.target lat5150-unified-api.service
Wants=network-online.target

[Service]
Type=oneshot
User=root
WorkingDirectory=/home/user/LAT5150DRVMIL/02-ai-engine

# Environment
Environment="PYTHONUNBUFFERED=1"
Environment="PYTHONPATH=/home/user/LAT5150DRVMIL/02-ai-engine:/home/user/LAT5150DRVMIL/03-web-interface"
Environment="HOME=/root"
Environment="PATH=/root/.local/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"

# Self-improvement configuration
Environment="AI_IMPROVEMENT_TARGET_SCORE=80.0"
Environment="AI_IMPROVEMENT_MAX_CYCLES=5"
Environment="AI_IMPROVEMENT_THRESHOLD=2.0"

# Run self-improvement session
ExecStart=/usr/bin/python3 /home/user/LAT5150DRVMIL/02-ai-engine/ai_self_improvement.py run

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=lat5150-self-improvement

# Resource limits
MemoryLimit=2G
CPUQuota=80%

# Timeout (max 30 minutes for full improvement session)
TimeoutStartSec=1800

# Security
PrivateTmp=true
ProtectSystem=strict
ReadWritePaths=/home/user/LAT5150DRVMIL/02-ai-engine/self_improvement_sessions
ReadWritePaths=/home/user/LAT5150DRVMIL/02-ai-engine/redteam_benchmark_data
ProtectHome=read-only
NoNewPrivileges=false

[Install]
WantedBy=multi-user.target
