# lat5150drvmil/deployment/systemd/lat5150-tactical.service
[Unit]
Description=LAT5150 DRVMIL Tactical Self-Coding AI Interface
Documentation=file:///home/user/LAT5150DRVMIL/00-documentation/
After=network.target network-online.target
Wants=network-online.target
ConditionPathExists=/home/user/LAT5150DRVMIL/03-web-interface/secured_self_coding_api.py

[Service]
Type=simple
User=root
WorkingDirectory=/home/user/LAT5150DRVMIL
Environment="PYTHONUNBUFFERED=1"
Environment="SECURITY_LEVEL=HIGH"
Environment="TACTICAL_PORT=5001"

# Start the secured tactical API
ExecStart=/usr/bin/python3 /home/user/LAT5150DRVMIL/03-web-interface/secured_self_coding_api.py \
    --workspace /home/user/LAT5150DRVMIL \
    --port 5001 \
    --security-level HIGH \
    --enable-rag \
    --enable-int8 \
    --enable-learning

# Restart policy
Restart=on-failure
RestartSec=10
StartLimitInterval=200
StartLimitBurst=5

# Resource limits
LimitNOFILE=65536
Nice=-10

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=lat5150-tactical

# Security hardening
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/home/user/LAT5150DRVMIL
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
